# deployments/systemd/mine-teleop-vehicle-agent.service
[Unit]
Description=Mine Teleop vehicle control agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/mine-teleop
EnvironmentFile=/etc/mine-teleop/mine-teleop.env
LogsDirectory=mine-teleop
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
ExecStartPre=/opt/mine-teleop/lib/ld-linux-x86-64.so.2 --library-path /opt/mine-teleop/lib:/opt/mine-teleop/lib/vendor/chassis:/opt/mine-teleop/lib/vendor/mvs:/opt/mine-teleop/lib/vendor/pylon /opt/mine-teleop/bin/mine-teleop vehicle-agent --config /etc/mine-teleop/vehicle-agent.yaml --preflight
ExecStartPre=/opt/mine-teleop/lib/ld-linux-x86-64.so.2 --library-path /opt/mine-teleop/lib:/opt/mine-teleop/lib/vendor/chassis:/opt/mine-teleop/lib/vendor/mvs:/opt/mine-teleop/lib/vendor/pylon /opt/mine-teleop/bin/mine-teleop vehicle-agent --config /etc/mine-teleop/vehicle-agent.yaml --adapter-status
ExecStart=/opt/mine-teleop/lib/ld-linux-x86-64.so.2 --library-path /opt/mine-teleop/lib:/opt/mine-teleop/lib/vendor/chassis:/opt/mine-teleop/lib/vendor/mvs:/opt/mine-teleop/lib/vendor/pylon /opt/mine-teleop/bin/mine-teleop vehicle-agent --config /etc/mine-teleop/vehicle-agent.yaml --teleop --service
Restart=on-failure
RestartSec=2
StandardOutput=append:/var/log/mine-teleop/vehicle-agent.log
StandardError=append:/var/log/mine-teleop/vehicle-agent.log

[Install]
WantedBy=multi-user.target
