# deployments/systemd/mine-teleop-vehicle-media-agent.service
[Unit]
Description=Mine Teleop vehicle media agent
After=network-online.target mine-teleop-vehicle-agent.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/mine-teleop
EnvironmentFile=/etc/mine-teleop/mine-teleop.env
LogsDirectory=mine-teleop
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
StateDirectory=mine-teleop/recordings
ExecStart=/opt/mine-teleop/mine-teleop vehicle-media-agent --config /etc/mine-teleop/vehicle-agent.yaml --driver-console-url ${MINE_TELEOP_DRIVER_CONSOLE_URL} --record --recording-root /var/lib/mine-teleop/recordings --service
Restart=on-failure
RestartSec=2
StandardOutput=append:/var/log/mine-teleop/vehicle-media-agent.log
StandardError=append:/var/log/mine-teleop/vehicle-media-agent.log

[Install]
WantedBy=multi-user.target
