# deployments/systemd/mine-teleop-turn-server.service
[Unit]
Description=Mine Teleop TURN Relay
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/mine-teleop
EnvironmentFile=/etc/mine-teleop/mine-teleop.env
LogsDirectory=mine-teleop
ExecStart=/usr/bin/turnserver -c /etc/mine-teleop/turnserver.conf --no-cli
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/mine-teleop/turn-server.log
StandardError=append:/var/log/mine-teleop/turn-server.err.log

[Install]
WantedBy=multi-user.target
